# board/raspberrypi/overlay/lib/systemd/system/wpa_supplicant.service
[Unit]
Description=WPA supplicant
Before=network.target
Wants=network.target
Wants=systemd-rfkill.service
After=systemd-rfkill.service

[Service]
Type=exec
ExecStart=/sbin/wpa_supplicant -i wlan0 -c /etc/wpa_supplicant_moosenet.conf

[Install]
WantedBy=multi-user.target
